# deploy/bobatea-backend.service
[Unit]
Description=Bobatea FastAPI backend
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/srv/bobatea/backend
EnvironmentFile=-/etc/bobatea/backend.env
ExecStart=/srv/bobatea/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
